# firmware/uci-defaults/99-default-theme
#!/bin/sh
# default theme: shadcn (the pipeline's only uci-defaults customization)
uci set luci.main.mediaurlbase='/luci-static/shadcn'
uci commit luci
exit 0
